# deploy/field-station.service
[Unit]
Description=Field Station — Claude Code visual config manager
After=network.target

[Service]
Type=simple

# REQUIRED: Set ExecStart to the absolute path of the field-station binary
ExecStart=/usr/local/bin/field-station

Restart=on-failure
RestartSec=5

# Optional: set auth token for password protection
# Environment=FIELD_STATION_TOKEN=your-secret-token

# Optional: enable secure cookies when behind an HTTPS proxy
# Environment=FIELD_STATION_SECURE_COOKIES=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=field-station

[Install]
WantedBy=default.target
