# deploy/field-station.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fieldstation.app</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>.output/server/index.mjs</string>
  </array>

  <!-- REQUIRED: Set this to the absolute path of your field-station directory -->
  <key>WorkingDirectory</key>
  <string>/path/to/field-station</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>~/Library/Logs/field-station/out.log</string>

  <key>StandardErrorPath</key>
  <string>~/Library/Logs/field-station/err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>
</dict>
</plist>
