# deploy/field-station.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fieldstation.app</string>

  <key>ProgramArguments</key>
  <array>
    <!-- REQUIRED: Set this to the absolute path of the field-station binary -->
    <string>/usr/local/bin/field-station</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>~/Library/Logs/field-station/out.log</string>

  <key>StandardErrorPath</key>
  <string>~/Library/Logs/field-station/err.log</string>

  <!-- Optional environment variables:
  <key>EnvironmentVariables</key>
  <dict>
    <key>FIELD_STATION_TOKEN</key>
    <string>your-secret-token</string>
    <key>FIELD_STATION_SECURE_COOKIES</key>
    <string>1</string>
  </dict>
  -->
</dict>
</plist>
